# products/virtual_tower/debian/mock_radio_virtual_tower.postinst
#!/bin/sh
set -e

if [ "$1" = "configure" ]
then
    # Start the service and enable it to start on boot.
    if [ -x "`which systemctl 2>/dev/null`" ]
    then
	for service in "tvsc-mock-communications" "tvsc-datetime" "tvsc-echo" "tvsc-hello" "tvsc-virtual-tower"
	do
	    systemctl start "${service}"
	    systemctl enable "${service}"
	done
    fi
fi
