# service/communications/debian/tvsc_communications_service.prerm
#!/bin/sh
set -e

if [ "$1" = "remove" ]; then
    if [ -x "`which systemctl 2>/dev/null`" ]; then
	systemctl stop tvsc-communications.service
	systemctl disable tvsc-communications.service
    fi
fi
