# service/hello/debian/tvsc_hello_service.prerm
#!/bin/sh
set -e

if [ "$1" = "remove" ]; then
    if [ -x "`which systemctl 2>/dev/null`" ]; then
	systemctl stop tvsc-hello.service
	systemctl disable tvsc-hello.service
    fi
fi
