# services/web_ui/debian/tvsc_web_ui_service.postinst
#!/bin/sh
set -e

if [ "$1" = "configure" ]; then
    # Start the service and enable it to start on boot.
    if [ -x "`which systemctl 2>/dev/null`" ]; then
	systemctl start tvsc-web-ui.service
	systemctl enable tvsc-web-ui.service
    fi
	    
fi
