# service/virtual_tower_ui/debian/tvsc_virtual_tower_ui_service.postinst
#!/bin/sh
set -e

if [ "$1" = "configure" ]; then
    # Start the service and enable it to start on boot.
    if [ -x "`which systemctl 2>/dev/null`" ]; then
	systemctl start tvsc-virtual-tower-ui.service
	systemctl enable tvsc-virtual-tower-ui.service
    fi
	    
fi
