# services/proxy/debian/tvsc_proxy_service.prerm
#!/bin/sh
set -e

if [ "$1" = "remove" ]; then
    if [ -x "`which systemctl 2>/dev/null`" ]; then
	systemctl stop tvsc-proxy.service
	systemctl disable tvsc-proxy.service
    fi
fi
