# service/proxy/debian/tvsc_proxy_service.prerm
#!/bin/sh
set -e

if [ "$1" = "remove" ]; then
    if [ -x "`which systemctl 2>/dev/null`" ]; then
	systemctl stop tvsc-proxy.service
	systemctl disable tvsc-proxy.service
	/sbin/iptables -t nat -D PREROUTING -p tcp --dport 80 -j REDIRECT --to-port 50050
	/sbin/iptables-save -f /etc/iptables/rules.v4
	/sbin/ip6tables-save -f /etc/iptables/rules.v6
    fi
fi
